# files/btrfs/root/etc/systemd/system/pine-first-boot.service
[Unit]
Description=Pine first boot script

[Service]
ExecStart=/usr/local/sbin/pine-first-boot
Type=oneshot

[Install]
WantedBy=multi-user.target
